# ansible/roles/pihole/templates/nebula-sync.service.j2
# Managed by Ansible — do not edit manually
[Unit]
Description=nebula-sync Pi-hole sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/nebula-sync.env
ExecStart=/usr/local/bin/nebula-sync