# ansible/roles/pihole/templates/nebula-sync.timer.j2
# Managed by Ansible — do not edit manually
[Unit]
Description=nebula-sync Pi-hole sync timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=7d

[Install]
WantedBy=timers.target